# config/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:nginx]
startretries=10
command=/usr/sbin/nginx

[program:uwsgi]
startretries=10
command=/usr/local/bin/uwsgi --ini /var/www/app/uwsgi.ini
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
